# config/systemd/user/notes.service
[Unit]
Description=sync notes

[Service]
Type=oneshot
ExecStart=/home/olex/org/upd.sh
